# meta-swi-sdx55/recipes-bsp/reboot-daemon/files/reboot-daemon.service
[Unit]
SourcePath=/sbin/reboot-daemon
Description= reboot-daemon service

[Service]
Restart=always
ExecStart=/sbin/reboot-daemon

[Install]
WantedBy=multi-user.target
